# deploy/raspberry-pi/energiemanager.service
[Unit]
Description=EnergyManager – Backend (API + Energiemanagement)
Documentation=https://github.com/daTobi1/Energiemanager
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=energiemanager
Group=energiemanager
WorkingDirectory=/home/energiemanager/energiemanager/backend
EnvironmentFile=/home/energiemanager/energiemanager/.env
ExecStartPre=docker compose -f /home/energiemanager/energiemanager/docker-compose.yml up -d db redis grafana
ExecStart=/home/energiemanager/energiemanager/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=10
TimeoutStartSec=120

# Sicherheit
NoNewPrivileges=true
ProtectHome=read-only
ReadWritePaths=/home/energiemanager/energiemanager

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=energiemanager

[Install]
WantedBy=multi-user.target
